# systemd/user/systemd.service
[Unit]
Description=Set GTK Theme based on the sun position

[Service]
ExecStart=%h/.local/bin/sun-sets-gtk-theme --dark-theme Adwaita-dark --light-theme Adwaita

[Install]
WantedBy=default.target
